# templates/irc-slack.service.erb
[Unit]
Description=IRC-to-Slack gateway
After=network.target
Documentation=<%= @soucre %>

[Service]
User=<%= @user %>
ExecStart=/usr/local/bin/irc-slack -L <%= @loglevel %> -H <%= @listen %> -p <%= @port %> -s <%= @server_name %>
Restart=always

[Install]
WantedBy=multi-user.target
